# dovecot-tmpfs-restore.service
[Unit]
Description=Restore Dovecot tmpfs directories from backup
Before=dovecot.service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/restore-dovecot-tmpfs.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target